# duplicity/files/duplicity.timer.jinja
[Unit]
Description=Duplicity timer

[Timer]
OnCalendar={{ on_calendar }}
Unit=duplicity.service

[Install]
WantedBy=timers.target
